# ext/templates/deb/prerm.erb
#!/bin/sh

set -e
if [ -x "/etc/init.d/dujour" ]; then
    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d dujour stop || exit $?
    else
        /etc/init.d/dujour stop || exit $?
    fi
fi

